# scripts/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable sensor.service
systemctl start sensor.service
