# scripts/sensor.service
[Unit]
Description=Sensor service management
After=mosquitto.service rethinkdb.service network-online.target
Wants=network-online.target
PartOf=rethinkdb.service network-online.target

[Service]
ExecStart=/usr/bin/sensorservice -c /etc/sensorservice/config.json
KillMode=process
PrivateTmp=true

[Install]
WantedBy=multi-user.target
